# service/memcached/run
#!/bin/sh
exec 2>&1
exec /usr/bin/memcached -u nobody  -m 150m


